# mazuadm-api.service
[Unit]
Description=MazuADM API
After=network.target postgresql.service

[Service]
Type=simple
Environment=DATABASE_URL=postgres://mazuadm:mazuadm@127.0.0.1:5433/mazuadm
ExecStart=/home/seadog007/MazuADM/target/release/mazuadm-api
Restart=on-failure

[Install]
WantedBy=multi-user.target
